# misc/master_updater.timer
[Unit]
Description = master_updater target

[Timer]
Unit = master_updater.service
OnCalendar = *:0/15

[Install]
WantedBy=timers.target
